# systemd/himawari8.service
# $HOME/.config/systemd/user/himawaripy.service

[Unit]
Description=Update desktop background with satellite imagery

[Service]
Type=oneshot
ExecStart=bash -c "himawari8 --auto-offset --save-battery"
